# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg

[program:php-fpm7.0]
command = /usr/sbin/php-fpm7.0

[program:nginx]
command=/usr/sbin/nginx
